# headscale/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/usr/bin/execlineb -S1
# shellcheck shell=bash
