# headscale/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Starting nginx reverse proxy..."
exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf
