# headscale/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "[INFO] Starting nginx reverse proxy..."

# Get the ingress entry from the HA Supervisor API or env var (for testing)
INGRESS_ENTRY="${INGRESS_ENTRY:-}"
if [ -z "${INGRESS_ENTRY}" ] && [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    INGRESS_ENTRY=$(curl -s -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
        http://supervisor/addons/self/info 2>/dev/null \
        | jq -r '.data.ingress_entry // empty')
fi

if [ -z "${INGRESS_ENTRY}" ]; then
    echo "[WARNING] Could not get ingress entry from Supervisor, using empty prefix"
fi

echo "[INFO] Ingress entry: ${INGRESS_ENTRY}"

# Wait for API key to be available (generated by headplane service)
API_KEY=""
if [ -f /data/headscale/api_key ]; then
    API_KEY=$(cat /data/headscale/api_key)
else
    echo "[INFO] Waiting for API key (up to 90s)..."
    for i in $(seq 1 90); do
        if [ -f /data/headscale/api_key ]; then
            API_KEY=$(cat /data/headscale/api_key)
            break
        fi
        sleep 1
    done
fi

if [ -n "${API_KEY}" ]; then
    echo "[INFO] API key found, auto-login enabled"
else
    echo "[WARNING] No API key found, auto-login disabled"
fi

# Replace placeholders in nginx config
# Use a delimiter that won't appear in the values (| instead of /)
sed -i "s|%%INGRESS_ENTRY%%|${INGRESS_ENTRY}|g" /etc/nginx/nginx.conf
sed -i "s|%%API_KEY%%|${API_KEY}|g" /etc/nginx/nginx.conf

exec nginx -c /etc/nginx/nginx.conf
