# headscale/rootfs/etc/s6-overlay/s6-rc.d/init-headscale/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
exec /usr/local/bin/generate-config.sh
