# headscale/rootfs/etc/s6-overlay/s6-rc.d/headscale/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
echo "[INFO] Headscale server stopped"
