# headscale/rootfs/etc/s6-overlay/s6-rc.d/headscale/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Headscale server stopped"
