# headscale/rootfs/etc/s6-overlay/s6-rc.d/headscale/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Starting Headscale server..."
exec /usr/local/bin/headscale serve --config /data/headscale/config.yaml
