# headscale/rootfs/etc/s6-overlay/s6-rc.d/headscale/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
echo "[INFO] Starting Headscale server..."
exec /usr/local/bin/headscale serve --config /data/headscale/config.yaml
