# headscale/rootfs/etc/s6-overlay/s6-rc.d/headplane/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
echo "[INFO] Headplane web UI stopped"
