# headscale/rootfs/etc/s6-overlay/s6-rc.d/headplane/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Headplane web UI stopped"
