# headscale/rootfs/etc/s6-overlay/s6-rc.d/headplane/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

_log() { echo "[INFO] $*"; }
_warn() { echo "[WARNING] $*"; }

# Read headscale port from options
if [ -n "${SUPERVISOR_TOKEN:-}" ] && command -v bashio::config &>/dev/null; then
    HEADSCALE_PORT=$(bashio::config 'headscale_port')
elif [ -f /data/options.json ]; then
    HEADSCALE_PORT=$(jq -r '.headscale_port' /data/options.json)
else
    HEADSCALE_PORT=8080
fi

_log "Starting Headplane web UI..."

# Wait for headscale to be ready
_log "Waiting for Headscale to be ready..."
for i in $(seq 1 30); do
    if curl -sf "http://127.0.0.1:${HEADSCALE_PORT}/health" > /dev/null 2>&1; then
        _log "Headscale is ready"
        break
    fi
    if [ "$i" -eq 30 ]; then
        _warn "Headscale did not become ready in time, starting Headplane anyway"
    fi
    sleep 1
done

# Generate API key on first run if needed
if [ ! -f /data/headscale/api_key ] && [ -f /usr/local/bin/generate-api-key.sh ]; then
    _log "Generating API key for Headplane..."
    /usr/local/bin/generate-api-key.sh
fi

# Validate the API key still works (may be stale after addon rebuild)
if [ -f /data/headscale/api_key ]; then
    API_KEY=$(cat /data/headscale/api_key)
    HTTP_CODE=$(curl -s -o /dev/null -w '%{http_code}' \
        -H "Authorization: Bearer ${API_KEY}" \
        "http://127.0.0.1:${HEADSCALE_PORT}/api/v1/apikey" 2>/dev/null)
    if [ "$HTTP_CODE" != "200" ]; then
        _warn "API key is invalid, regenerating..."
        rm -f /data/headscale/api_key
        /usr/local/bin/generate-api-key.sh
    fi
fi

# Show API key in logs so user can find it in the Log tab
if [ -f /data/headscale/api_key ]; then
    _log "========================================="
    _log "Headplane API Key (paste into the login page):"
    _log "$(cat /data/headscale/api_key)"
    _log "========================================="
fi

export NODE_ENV=production
export HEADPLANE_CONFIG=/etc/headplane/config.yaml
cd /opt/headplane
exec node /opt/headplane/build/server/index.js
